# snova/config/templates/systemd/sparrow-snova-sparrow-default-worker.service
[Unit]
Description="{{ snova_name }}-sparrow-default-worker %I"
PartOf={{ snova_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ snova_cmd }} worker --queue default
StandardOutput=file:{{ snova_dir }}/logs/worker.log
StandardError=file:{{ snova_dir }}/logs/worker.error.log
WorkingDirectory={{ snova_dir }}
