# snova/config/templates/systemd/sparrow-snova-web.target
[Unit]
After=network.target
Wants={{ snova_name }}-frappe-web.service {{ snova_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
