# snova/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ snova_name }}-sparrow-web]
command={{ snova_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} --max-requests {{ gunicorn_max_requests }} --max-requests-jitter {{ gunicorn_max_requests_jitter }} -t {{ http_timeout }} sparrow.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/web.log
stderr_logfile={{ snova_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ snova_name }}-sparrow-schedule]
command={{ snova_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/schedule.log
stderr_logfile={{ snova_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ snova_dir }}

{% if not multi_queue_consumption %}
[program:{{ snova_name }}-sparrow-default-worker]
command={{ snova_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/worker.log
stderr_logfile={{ snova_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ snova_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endif %}

[program:{{ snova_name }}-sparrow-short-worker]
command={{ snova_cmd }} worker --queue short{{',default' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/worker.log
stderr_logfile={{ snova_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ snova_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ snova_name }}-sparrow-long-worker]
command={{ snova_cmd }} worker --queue long{{',default,short' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/worker.log
stderr_logfile={{ snova_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ snova_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% for worker_name, worker_details in workers.items() %}
[program:{{ snova_name }}-sparrow-{{ worker_name }}-worker]
command={{ snova_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/worker.log
stderr_logfile={{ snova_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ snova_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endfor %}


{% if not skip_redis %}
[program:{{ snova_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/redis-cache.log
stderr_logfile={{ snova_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ snova_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/redis-queue.log
stderr_logfile={{ snova_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if node %}
[program:{{ snova_name }}-node-socketio]
command={{ node }} {{ snova_dir }}/apps/sparrow/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ snova_dir }}/logs/node-socketio.log
stderr_logfile={{ snova_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ snova_dir }}
{% endif %}

[group:{{ snova_name }}-web]
programs={{ snova_name }}-sparrow-web {%- if node -%} ,{{ snova_name }}-node-socketio {%- endif%}


{% if multi_queue_consumption %}

[group:{{ snova_name }}-workers]
programs={{ snova_name }}-sparrow-schedule,{{ snova_name }}-sparrow-short-worker,{{ snova_name }}-sparrow-long-worker{%- for worker_name in workers -%},{{ snova_name }}-sparrow-{{ worker_name }}-worker{%- endfor %}

{% else %}

[group:{{ snova_name }}-workers]
programs={{ snova_name }}-sparrow-schedule,{{ snova_name }}-sparrow-default-worker,{{ snova_name }}-sparrow-short-worker,{{ snova_name }}-sparrow-long-worker{%- for worker_name in workers -%},{{ snova_name }}-sparrow-{{ worker_name }}-worker{%- endfor %}

{% endif %}

{% if not skip_redis %}
[group:{{ snova_name }}-redis]
programs={{ snova_name }}-redis-cache,{{ snova_name }}-redis-queue
{% endif %}
